# crontab
PATH=/usr/bin:/bin:/home/[USER]/speedtest
42 * * * *      monitor.sh 2>&1 > /dev/null
